# docker/neutron/rabbitmq.conf
[program:rabbitmq-server]
command=/usr/sbin/rabbitmq-server
redirect_stderr=true
autostart=true
